# resources/gru.service
[Unit]
Description=Gru
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/gru serve
